# src/apps/halcsd/init/systemd/etc/udev/rules.d/60-udev_fpga.rules.in
# Permissions for the fpga devices

KERNEL=="fpga*", ACTION=="add", OWNER="root", GROUP="users", MODE="660", PROGRAM="@halcsd_HALCS_INSTALL_FULL_DATADIR@/generate-dev-idx.sh %k", SYMLINK+="fpga/%c{1}", RUN+="/usr/bin/systemctl --no-block start halcs-ioc@%c{2}.target halcs-ioc@%c{3}.target"

KERNEL=="fpga*", ACTION=="remove", PROGRAM="@halcsd_HALCS_INSTALL_FULL_DATADIR@/generate-dev-idx.sh %k", RUN+="/usr/bin/systemctl --no-block stop halcs@%c{2}.target halcs@%c{3}.target"
